# contrib/systemd/noterious.service
[Unit]
Description=Noterious personal knowledge base
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/usr/bin/noterious
Restart=on-failure
RestartSec=2s

# Packaged installs can enable this user unit directly.
# Source installs can copy it to ~/.config/systemd/user/noterious.service
# and adjust ExecStart plus the paths below as needed.
# All environment variables currently supported by the server are listed here.
# Leave the bootstrap variables empty for the normal first-run web setup flow.
# Only set NOTERIOUS_AUTH_BOOTSTRAP_PASSWORD for unattended/headless provisioning.
# NOTERIOUS_AUTH_BOOTSTRAP_PASSWORD_FILE can be used instead when a secret file is available.
# NOTERIOUS_AUTH_BOOTSTRAP_USERNAME is only used when the bootstrap password is set.
Environment="NOTERIOUS_LISTEN_ADDR=:3000"
Environment="NOTERIOUS_VAULT_PATH=%h/Documents/noterious"
Environment="NOTERIOUS_DATA_DIR=%h/.local/state/noterious"
Environment="NOTERIOUS_WATCH_INTERVAL=2s"
Environment="NOTERIOUS_NTFY_INTERVAL=1m"
Environment="NOTERIOUS_AUTH_COOKIE_NAME=noterious_session"
Environment="NOTERIOUS_AUTH_SESSION_TTL=720h"
Environment="NOTERIOUS_AUTH_BOOTSTRAP_USERNAME="
Environment="NOTERIOUS_AUTH_BOOTSTRAP_PASSWORD="
Environment="NOTERIOUS_AUTH_BOOTSTRAP_PASSWORD_FILE="

[Install]
WantedBy=default.target
